# container/cron.d/mcwzh-croe-diff
0 0 * * * dinitctl trigger mcwzh-croe-diff
