# container/cron.d/mcwzh-walker
0 0 */7 * * /dist/sbin/dinitctl start mcwzh-walker
