# container/cron.d/mcwzh-ia-archiver
0 30 3 * * /dist/sbin/dinitctl start mcwzh-ia-archiver
